# files/etc/init/monit.conf
# monit -*- sh -*-
#
# Make sure monit is respawn if it dies.
#

start on runlevel [2345]
stop on runlevel [06]

exec /usr/sbin/monit -Ic /etc/monitrc
respawn
